# devops/familysync.service
[Unit]
Description=FamilySync Backend Server
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/familysync
EnvironmentFile=/opt/familysync/.env
ExecStart=/opt/familysync/familysync-server
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=familysync

# Resource limits — sensible for 2 GB RAM Armbian
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target
